# docker/builder/profile.d/default
-DENABLE_CLICKHOUSE_BENCHMARK=ON
-DENABLE_CURL=ON
-DENABLE_S3=ON
-DENABLE_SSL=ON
-DENABLE_LICENSE_CHECK=OFF
